# agentic/pi/extensions/telegram-bridge/pi-telegram-bridge.service
[Unit]
Description=Pi Telegram Bridge
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/dotfiles/agentic/pi/extensions/telegram-bridge
ExecStart=%h/.local/bin/uv run src/scripts/bridge.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
